# init/cron/crontab
#crontab
0-59/10 * * * * unxdev01 *default/cron_unfinished_actions
59 23 * * * unxdev01 *default/cron_stats
55 23 * * * unxdev01 *default/refresh_obsolescence
0 1 * * * unxdev01 *default/viz_cron_cleanup
1 5 * * * unxdev01 *alerts/cron_alerts_daily
1 * * * * unxdev01 *alerts/cron_alerts_hourly
